# systemd/aios-helper.service
[Unit]
Description=AIOS Privileged Helper
After=local-fs.target
Before=aios-watchdog.service aios-agent.service

[Service]
Type=simple
ExecStart=/usr/local/bin/aios-helper
Restart=no
User=root
Group=root
RuntimeDirectory=aios
RuntimeDirectoryMode=0750
StandardOutput=journal
StandardError=journal
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
