# systemd/aios-helper.service
[Unit]
Description=AIOS Privileged Helper
After=network.target
Before=aios-watchdog.service

[Service]
Type=simple
ExecStart=/usr/local/bin/aios-helper
Restart=no
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target